# conf/supervisor-hbase.ini
[program:hbase]
command=/data/start_hbase.sh
stopwaitsecs=60
