# conf/supervisor-opentsdb.ini
[program:opentsdb]
command="COMPRESSION=gzip /data/start_opentsdb.sh"
stopwaitsecs=60
