# conf/supervisor-opentsdb.ini
[program:opentsdb]
#environment="COMPRESSION=gzip"
command="/data/start-opentsdb.sh"
stopwaitsecs=60
