# conf/supervisor-hbase.ini
[program:hbase]
command=/data/start-hbase.sh
stopwaitsecs=60
