# conf/supervisor-opentsdb.ini
[program:opentsdb]
command=/tsdb/start_opentsdb.sh
stopwaitsecs=60
