# conf/supervisor-scollector.ini
[program:scollector]
command=/scollector/scollector -h localhost:8070
priority=30
