# conf/supervisor-bosun.ini
[program:bosun]
command=/bosun/bosun -c /data/bosun.conf
priority=30
